# profile.d/clang-format.sh
alias clang-format-deep="find . -regex '^.*\.\(c\|cpp\|h\|hpp\)$' | xargs clang-format -i"
