# profile.d/flatpak.sh
alias code="flatpak run com.visualstudio.code"
alias zed="flatpak run dev.zed.Zed"
alias typora="flatpak run io.typora.Typora"
alias calc="flatpak run org.kde.kcalc"
