# profile.d/flatpak.sh
alias code="flatpak run com.visualstudio.code"
alias zed="flatpak run dev.zed.Zed"
